# default/profile.d/99-config-envvars.sh
export CONFIG_PATH="$HOME/.local/share/arch-config"
export PATH="$CONFIG_PATH/bin:$PATH"

# Enable MangoHud by default
export MANGOHUD=1
